# clamav-scan-daily.service
[Service]
Type=oneshot
User=root

# Be gentle to other workloads
Nice=15
IOSchedulingClass=idle

# Treat "virus found" (exit 1) as success
SuccessExitStatus=0 1

# Prevent runaway jobs
RuntimeMaxSec=1h

# Ensure log directory exists
ExecStartPre=/usr/bin/install -o root -g root -m 750 -d /var/log/clamav

# FAST scan: daemon-based, parallel directory scanning, FD passing for perms
ExecStart=/usr/bin/clamdscan --fdpass --multiscan --infected \
  --log=/var/log/clamav/daily-scan.log \
  /home /var/www /opt

# Hardening (still allows scanning)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target

