# clamav-scan-weekly.service
[Unit]
Description=ClamAV Weekly Full System Scan (clamscan)
After=network-online.target clamav-freshclam.service
Wants=clamav-freshclam.service

[Service]
Type=oneshot
# Run as root for broad filesystem access and fewer permission errors.
User=clamav

# Do not mark the unit failed when malware is found (exit 1).
SuccessExitStatus=0 1

# Keep scans from running forever; weekly full scans can be long.
RuntimeMaxSec=12h

# Be polite to the system.
Nice=10
IOSchedulingClass=best-effort

# Protect logs that may contain sensitive paths/filenames.
UMask=0077

# Make sure log dir exists and definitions are current.
ExecStartPre=/usr/bin/install -d -m 0755 /var/log/clamav
ExecStartPre=/usr/bin/freshclam -q

# The scan itself. Use anchored regex for system dirs. Limit to a single filesystem if desired.
ExecStart=/usr/bin/clamscan -r / \
  --infected \
  --log=/var/log/clamav/weekly-full-scan.log \
  --exclude-dir='^/proc' \
  --exclude-dir='^/sys' \
  --exclude-dir='^/run' \
  --exclude-dir='^/dev' \
  --exclude-dir='^/var/lib/flatpak' \
  --cross-fs=no

# Reasonable hardening that still allows read access to the whole FS.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHome=no
ProtectSystem=no
